# ex04/82-nsikhosa-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/lil_penguin/ex04/main.ko"
